# deploy/sirel-bots.service
[Unit]
Description=Sirel Telegram + Max bots (start.py)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/Sirel_Concierge
Environment=PATH=/root/Sirel_Concierge/.venv/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
ExecStart=/root/Sirel_Concierge/.venv/bin/python3 /root/Sirel_Concierge/start.py
Restart=always
RestartSec=20
StartLimitIntervalSec=400
StartLimitBurst=8
KillSignal=SIGTERM
TimeoutStopSec=35

[Install]
WantedBy=multi-user.target
